# 11f38616-31a9-4359-9005-75103cf8e6d9_openkernel/files.list
init.lua
sbin/
sbin/init.lua
sbin/openrc.lua
lib/
lib/tokenize.lua
lib/colors.lua
lib/motd/
lib/motd/motdlist.lua
lib/require/
lib/require/serialize.lua
lib/services/
lib/services/logd.lua
etc/
etc/hostname.lua
etc/openrc.cfg
etc/init.d/
etc/init.d/00_base.lua
etc/init.d/10_event.lua
etc/init.d/20_read.lua
etc/init.d/30_errors.lua
etc/init.d/40_tables.lua
etc/init.d/50_modules.lua
etc/init.d/60_networking.lua
etc/init.d/70_filesystem.lua
etc/init.d/80_multithreading.lua
boot/
boot/kernel.lua
boot/libs/
boot/libs/filesystem.lua
boot/libs/term.lua
bin/
bin/eeprom.lua
bin/label.lua
bin/shutdown.lua
bin/sh.lua
bin/rm.lua
bin/pwd.lua
bin/motd.lua
bin/cp.lua
bin/mv.lua
bin/mkdir.lua
bin/ls.lua
bin/free.lua
bin/exit.lua
bin/echo.lua
bin/clear.lua
bin/cd.lua
bin/cat.lua
bin/umount.lua
bin/mount.lua
bin/hostname.lua
bin/get.lua
bin/pastebin.lua
bin/wget.lua
usr/bin/
usr/bin/lua.lua
var/log/
